# home-server/systemd/system/synchronize-krathalan-repo.service
[Unit]
Description=Synchronize krack pkgcache with reverse proxy pacman cache

[Service]
Type=oneshot
ExecStart=/usr/bin/synchronize-krathalan-repo.sh
