# base/systemd/system/systemd-logind.service.d/hidepid.conf
# https://wiki.archlinux.org/title/Security#hidepid
[Service]
SupplementaryGroups=proc
