# server/systemd/system/krack-receive.service.d/protect-system-strict.conf
[Service]
ProtectSystem=strict
ReadWritePaths=/var/www/builds/x86_64/
