# server/systemd/system/backup-email.timer
[Unit]
Description=Weekly backup of /etc and /home

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
